# dnsmasq-manager.service
[Unit]
After=dnsmasq.service

[Service]
ExecStart=/usr/bin/dnsmasq-manager
Restart=on-failure

[Install]
WantedBy=multi-user.target
